# zakibot.service
[Unit]
Description=Zakibot Service - Manages CameraFeed and RobotMain
After=network.target

[Service]
# Use absolute paths; ensure that 'python3' points to the correct interpreter
ExecStart=/usr/bin/python3 /home/amjpi/Zakibot/start_bot.py

# Set the working directory if necessary
WorkingDirectory=/home/amjpi/

# Environment variables can be set here if required
# Uncomment and adjust the following line if using a virtual environment
# Environment=PATH=/home/amjpi/venv/bin:/usr/bin:/bin

# Redirect stdout and stderr to the journal
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=always
RestartSec=5

# Run as a specific user and group
User=amjpi
Group=amjpi

# Ensure that child processes are killed when the service stops
KillMode=process
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target

